# superivsord.conf
[program:jxzyk]
command=/home/deploy/jxzyk/venv/bin/uwsgi --ini /home/deploy/jxzyk/uwsgi.ini