# superivsord.conf
[program:jxzyk]
command=/home/ubuntu/pywebsites/jxzyk/venv/bin/uwsgi --ini /home/ubuntu/pywebsites/jxzyk/uwsgi.ini